# service/udiskie/run
#!/bin/sh
exec udiskie --automount --notify --no-tray
