# systemd/user/notmuch.timer
[Unit]
Description=Notmuch timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=5m

[Install]
WantedBy=default.target
